# deploy/systemd/auto-articles.service.example
# Copy to /etc/systemd/system/auto-articles.service
# sudo systemctl daemon-reload && sudo systemctl enable --now auto-articles
#
# Adjust paths and User= to match your VM (e.g. /var/www/auto-articles).

[Unit]
Description=Auto Articles (legacy Flask/gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/auto-articles
Environment=HOME=/var/www/auto-articles
Environment=FLASK_ENV=production
EnvironmentFile=/var/www/auto-articles/.env
ExecStart=/var/www/auto-articles/.venv/bin/gunicorn wsgi:app --bind 127.0.0.1:8000 --workers 2 --threads 4 --timeout 180
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

# NOTE:
# This unit is kept only as a reference for the legacy Flask app.
# For the current system, use:
# - deploy/systemd/auto-articles-backend.service.example (FastAPI)
# - deploy/systemd/auto-articles-frontend.service.example (Next.js)
